# .config/systemd/user/update-cloud-images.service
[Unit]
Description=Update cloud images in ~/images/

[Service]
Type=exec
ExecStart=/home/jamesvaughn/bin/update --images
