# .config/systemd/user/update-vagrant-boxes.service
[Unit]
Description=Update Vagrant boxes

[Service]
Type=exec
ExecStart=/home/jamesvaughn/bin/update --vagrant
