# .config/autostart/gnome-configure.desktop
[Desktop Entry]
Exec=/home/vagrant/.config/autostart/gnome-configure.sh
Type=Application
